# root/etc/services.d/ttrss-daemon/run
#!/bin/sh

while true; do
  cd /var/www/ttrss
  php -f /var/www/ttrss/update_daemon2.php
  sleep 5m
done
