# service-php5-fpm.conf
[program:php5-fpm]
command=/usr/sbin/php5-fpm
user=root
redirect_stderr=true
