# service-ttrss-update.conf
[program:ttrss-update]
command=/srv/update-ttrss.sh
user=root
redirect_stderr=true
