# service-ttrss-update.conf
[program:ttrss-update]
command=/srv/update-ttrss.sh --wait-exit 24h
user=root
autorestart=true
exitcodes=0
redirect_stderr=true
