# root/etc/services.d/ttrss-updater/run
#!/bin/sh

while true; do
  /srv/update-ttrss.sh
  sleep 24h
done
